# zabbix-kiosk.service
[Unit]
Description=Zabbix Kiosk Monitoring System
After=network.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/opt/zabbix-kiosk
Environment="PATH=/opt/zabbix-kiosk/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/zabbix-kiosk"
Environment="DATABASE_URL=sqlite:///./data/zabbix_kiosk.db"
Environment="SECRET_KEY=change-this-in-production-very-secret-key-min-32-chars"
Environment="PORT=8000"
ExecStart=/opt/zabbix-kiosk/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --app-dir backend
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target